# config/crontab
# send the cleanup signal to the process every 10 minutes
# if the presentation is not connected to by any users, remove it from memory
*/15 * * * * path/to/cleanup.sh
